# deploy/okk-meetings.cron
# Atamura OKK — ОП meeting pipeline every 30 minutes (cron alternative to systemd).
# Install:  crontab -u okk deploy/okk-meetings.cron
*/30 * * * * cd /opt/AtamuraOKK && .venv/bin/python -m AtamuraOKK.scoring.meetings run --limit 50 >> /var/log/okk-meetings.log 2>&1
